# proxpanel-tunnel.service
[Unit]
Description=ProxPanel Remote Support Tunnel
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/proxpanel/tunnel-service.sh
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
